# configs/supervisor.conf.dist
[program:main]
command=xinit /opt/xbterminal/xbterminal/main.py

[group:xbterminal]
programs=main
